# supervisor.conf
[program:ldev-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/vagrant/projects/ldev/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=vagrant
numprocs=8
redirect_stderr=true
stdout_logfile=/home/vagrant/projects/ldev/supervisor.log